# 4.2/supervisord/4-websocket.conf
[program:websocket]
command=/var/www/orocommerce/bin/console gos:websocket:server --env=prod
numprocs=1
autostart=false
autorestart=true
startsecs=0
user=www-data
redirect_stderr=true
